# supervisor.conf.d/adb.conf
[program:adb_daemon]
command=adb start-server && adb status-window 
